# discord-bot.service
[Unit]
Description=Discord Channel Reset Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/discord-bot
Environment=PATH=/home/YOUR_USERNAME/discord-bot/venv/bin
ExecStart=/home/YOUR_USERNAME/discord-bot/venv/bin/python3 bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
